# src/scripts/sensor-agent.service
[Unit]
Description=Smartfarm Sensor Agent (CPU/Memory Monitoring)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/smartfarm-sensor-agent
Environment="PATH=/usr/local/bin:/usr/bin"
ExecStart=/usr/bin/python3 sensor_agent.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sensor-agent

[Install]
WantedBy=multi-user.target
